# src/avrecd
#!/sbin/runscript
# Copyright 1999-2005 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

# dvrec

opts="${opts} reload"

depend() {
	use alsasound
	after modules
}

config=/etc/avrecord/avrecord.conf

start() {
	sleep 5
	ebegin "Starting avrecd"
	MALLOC_CHECK_=1 /usr/sbin/avrecord -d -c $config >> /etc/avrecord/avrecord.log 2>&1
	eend $? "Failed to start avrecd"
}

stop() {
	ebegin "Stopping avrecd"
	pkill avrecord
	eend $? "Failed to stop avrecd"
}

reload() {
	ebegin "Reloading avrecd configuration"
	kill -HUP $(pidof avrecord) > /dev/null 2>&1
	eend $? "Failed to reload avrecd configuration"
}
